# centos6/mariadb/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=info

[program:sshd]
command=/usr/sbin/sshd -D
