# centos7/build/apache/php/5.6/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=info

[program:apache]
command=/usr/sbin/httpd
stdout_events_enabled=true
stderr_events_enabled=true
autostart = true
autorestart = true
stopasgroup = true
