# centos7/build/wordpress/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=info

[program:php-fpm]
command=/opt/remi/php70/root/usr/sbin/php-fpm 
stdout_events_enabled=true
stderr_events_enabled=true
autostart = true
autorestart = true
stopasgroup = true

